# roles/settings/templates/etc/profile.d/lujeni.sh.j2
setxkbmap us intl
xrandr -s 1920x1080
xhost local:root
